# systemd/logsentinel.service
[Unit]
Description=logsentinel web defender

[Service]
Type=oneshot
ExecStart=/usr/local/bin/logsentinel run
