# systemd/logsentinel.service
[Unit]
Description=webguard web defender

[Service]
Type=oneshot
ExecStart=/usr/local/bin/webguard run
